# monitoring/systemd/promtail.service
[Unit]
Description=Promtail log shipper
After=network-online.target loki.service
Wants=network-online.target

[Service]
Type=simple
User=promtail
Group=promtail
ExecStart=/usr/local/bin/promtail -config.file=/etc/promtail/promtail-config.yml
Restart=always
RestartSec=5
CPUQuota=100%
MemoryMax=512M
TasksMax=50
OOMScoreAdjust=-500
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /tmp/positions.yaml
PrivateTmp=true

[Install]
WantedBy=multi-user.target
# VERIFIED: Promtail systemd unit with read access to /var/log.
